# deploy/family-financials.service.template
[Unit]
Description=Family Financials live site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env npm run live
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
