# deploy/plaid-webhook.service.template
[Unit]
Description=Family Financials Plaid webhook listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
EnvironmentFile=__WORKDIR__/.env.plaid-webhook
ExecStart=/usr/bin/env npm run plaid:webhook
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
